# configs/supervisor.d/caddy.ini
[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
user=root
autostart=true
autorestart=true