# configs/supervisor.d/crond.ini
[program:crond]
command=crond -f
user=root
autostart=true
autorestart=true