# configs/supervisor.d/php-fpm.ini
[program:php-fpm]
command=php-fpm
user=root
autostart=true
autorestart=true